# binwalk-2.0.1/_openwrt-ar71xx-generic-gl-inet-v1-squashfs-sysupgrade.bin.extracted/squashfs-root/etc/uci-defaults/90-lime-banner
#!/bin/sh

cat <<EOF >> /etc/banner
  ___   __ __                   _______             __
 |   |_|__|  |--.----.-----. _ |   |   |-----.-----|  |--.
 |     |  |  _  |   _|  -__||_||       |  -__|__ --|     |
 |_____|__|_____|__| |_____|   |__|_|__|_____|_____|__|__|

 ---------------------------------------------------------
 LiMe 14.08 bigbang (develop rev.3ee0272 20150120_1709)
 ---------------------------------------------------------
 http://libre-mesh.org
 ---------------------------------------------------------
EOF

exit 0
