# binwalk-2.0.1/_openwrt-ar71xx-generic-gl-inet-v1-squashfs-sysupgrade.bin.extracted/squashfs-root/usr/lib/opkg/info/alfred.list
/etc/alfred/bat-hosts.lua
/usr/sbin/alfred
/etc/init.d/alfred
/etc/config/alfred
/usr/sbin/batadv-vis
