# binwalk-2.0.1/_openwrt-ar71xx-generic-gl-inet-v1-squashfs-sysupgrade.bin.extracted/squashfs-root/usr/lib/opkg/info/dnsmasq.list
/etc/hotplug.d/iface/25-dnsmasq
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
