# binwalk-2.0.1/_openwrt-ar71xx-generic-gl-inet-v1-squashfs-sysupgrade.bin.extracted/squashfs-root/usr/lib/opkg/info/luci-app-splash.list
/usr/lib/lua/luci/i18n/splash.pt-br.lmo
/etc/config/luci_splash
/usr/lib/lua/luci/i18n/splash.it.lmo
/www/luci/splash/index.html
/usr/lib/lua/luci/i18n/splash.pl.lmo
/usr/lib/lua/luci/i18n/splash.pt.lmo
/usr/lib/lua/luci/view/splash_splash/splash.htm
/usr/lib/lua/luci/i18n/splash.es.lmo
/usr/lib/lua/luci/i18n/splash.fr.lmo
/lib/upgrade/keep.d/luci-app-splash
/usr/lib/lua/luci/view/admin_status/splash.htm
/etc/init.d/luci_splash
/usr/lib/lua/luci/i18n/splash.de.lmo
/usr/lib/lua/luci/view/splash/splash.htm
/etc/config/luci_splash_leases
/usr/lib/lua/luci/i18n/splash.el.lmo
/usr/sbin/luci-splash
/usr/lib/lua/luci/i18n/splash.zh-cn.lmo
/usr/lib/lua/luci/i18n/splash.ru.lmo
/usr/lib/lua/luci/i18n/splash.ca.lmo
/usr/lib/lua/luci/controller/splash/splash.lua
/usr/lib/lua/luci/view/splash/blocked.htm
/usr/lib/lua/luci/model/cbi/splash/splash.lua
/www/cgi-bin/splash/splash.sh
/etc/hotplug.d/iface/25-splash
/usr/lib/lua/luci/view/splash_splash/index.htm
/usr/lib/lua/luci/i18n/splash.cs.lmo
/usr/lib/lua/luci/model/cbi/splash/splashtext.lua
/usr/lib/lua/luci/i18n/splash.he.lmo
