# binwalk-2.0.1/_openwrt-ar71xx-generic-gl-inet-v1-squashfs-sysupgrade.bin.extracted/squashfs-root/usr/lib/opkg/info/watchping.list
/etc/watchping/wan-ok.d/README
/usr/bin/watchping
/etc/watchping/wan-fail.d/README
/etc/init.d/watchping
/etc/uci-defaults/90_watchping
